# rpms/nas_img-bin-1.5.5.1-1.t1.i686/etc/cron.daily/naswatch
#!/bin/sh
TMP_FILE_LIST=`find /tmp/ -type f`

echo "${TMP_FILE_LIST}" | while read FILE;do
	touch -c $FILE
done
